# systemd/.config/systemd/user/setbg.service
[Unit]
Description=Set Backround

[Service]
Type=simple
StandardOutput=journal
ExecStart=%h/.local/bin/setbg

[Install]
WantedBy=default.target
